# server_config/supervisor/conf.d/gunicorn_channels.conf
[program:channels-gunicorn]
# TCP socket used by Nginx backend upstream
socket=tcp://localhost:8000

# Directory where your site's project files are located
directory=/home/pi/myProjects/channel

# Each process needs to have a separate socket file, so we use process_num
# Make sure to update "mysite.asgi" to match your project name
# command=/home/pi/.virtualenvs/channels35/bin/gunicorn --workers 4 --bind unix:/var/run/gunicorn/gunicorn-channels.sock skaterhome.wsgi:application
command=/home/pi/.virtualenvs/channels35/bin/gunicorn --workers 4 --bind 127.0.0.1:8080 skaterhome.wsgi:application

#process_name=gunicorn-channels%(process_num)s
#numprocs=4

# Automatically start and recover processes
autostart=true
autorestart=true

# Choose where you want your log to go
stdout_logfile=/var/log/gunicorn-channels.log
redirect_stderr=true
